# ons-business-dev/provisioning/roles/apache-cassandra/templates/etc/profile.d/apache-cassandra.sh
export PATH=$PATH:{{apache_cassandra_home}}/bin
