# ons-automation/roles/apache-hadoop/templates/etc/profile.d/apache-hadoop.sh
export PATH=$PATH:{{apache_hadoop_home}}/bin:{{apache_hadoop_home}}/sbin
