# ons-automation/roles/apache-spark/templates/etc/profile.d/apache-spark.sh
export PATH=$PATH:{{apache_spark_home}}/bin
